# examples/template/etc/init/bloonix-webapp.service.in
[Unit]
Description=Bloonix Web
After=network.target

[Service]
Type=forking
PIDFile=@@RUNDIR@@/bloonix/bloonix-@@NAME@@.pid
ExecStartPre=@@USRLIBDIR@@/bloonix/bin/bloonix-pre-start @@CACHEDIR@@/bloonix @@LIBDIR@@/bloonix @@LOGDIR@@/bloonix @@RUNDIR@@/bloonix
ExecStart=@@SRVDIR@@/bloonix/@@NAME@@/scripts/bloonix-@@NAME@@ --pid-file @@RUNDIR@@/bloonix/bloonix-@@NAME@@.pid --config-file @@CONFDIR@@/bloonix/@@NAME@@/main.conf
ExecReload=/bin/kill -1 $MAINPID
ExecStop=/bin/kill -15 $MAINPID
PrivateTmp=true
TimeoutStartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
